# whorang/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: WhoRang - Nginx Service
# ==============================================================================

bashio::log.info "Starting Nginx service..."

# FIXED: Run nginx as the correct user (abc) with proper configuration
exec s6-setuidgid abc nginx -g "daemon off;"
