# whorang/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bash

# Start nginx in the foreground for s6-overlay supervision
exec /usr/sbin/nginx -g 'daemon off;'
    sed -i '/allow 172\.30\.32\.2;/,/# Note: In standalone mode, run\.sh will replace this with/c\
    # Standalone mode - allow all access\
    allow all;' /etc/nginx/conf.d/default.conf 2>/dev/null || true
fi

# Validate nginx configuration
if ! nginx -t; then
    log_error "Nginx configuration test failed"
    exit 1
fi

# Start nginx in foreground mode
exec nginx -g "error_log /dev/stderr info; daemon off; master_process off;"
