# whorang/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bash

# Exit immediately if a command exits with a non-zero status.
set -e

# In standalone mode, we need to relax the IP restrictions.
# The WHORANG_ADDON_MODE variable is set by the 01-whorang-init script.
if [ "${WHORANG_ADDON_MODE}" = "false" ]; then
    echo "[nginx-run] Standalone mode detected. Allowing all IPs."
    # This sed command replaces the block of 'allow' rules with 'allow all;'
    # It targets backend.conf, which is the correct file.
    sed -i '/allow 172\.30\.32\.2;/,/allow ::1;/c\    allow all;' /etc/nginx/conf.d/backend.conf
fi

# Always validate the final nginx configuration before starting.
echo "[nginx-run] Validating nginx configuration..."
if ! /usr/sbin/nginx -t; then
    echo "[nginx-run] Nginx configuration test failed. Exiting."
    exit 1
fi
echo "[nginx-run] Nginx configuration is valid."

# Start nginx in the foreground for s6-overlay supervision.
# This is the proper way to run nginx under s6.
# We pass necessary directives to prevent daemonization and log to stderr.
echo "[nginx-run] Starting nginx..."
exec /usr/sbin/nginx -g "daemon off;"
