# whorang/rootfs/etc/services.d/nodejs/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: WhoRang - Node.js Service
# ==============================================================================

bashio::log.info "Starting WhoRang backend service..."

# Change to the app directory
cd /app || exit 1

# FIXED: Use standard HA addon user (abc) instead of whorun
# Execute the Node.js server with proper user and environment
exec s6-setuidgid abc node server.js
